# cdm_bot.service
[Unit]
Description=CDM Discount Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/cdm_bot
EnvironmentFile=/opt/cdm_bot/.env
ExecStart=/opt/cdm_bot/venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
